# roles/deps/files/alet.crontab
$ANSIBLE_VAULT;1.1;AES256
66353562306231306238366236313961346336333935316161643631323439353535616531623666
3431336265313633633037653566353332323439313930380a336534306366326339303235633262
64623861306536313738346230396264336264353433356231313464373838316133366133613134
3833613533333763660a643965373563373638656164343761353363333562656539353336613238
62656234373535363032333239656538356435336137663637343063633435626230663239333630
37636365623330313561386330346563626165333532343339353933663733613530356234373565
61326263653564646462656631643039666133666436376466643466316364313739383038303036
37376434303136343165323866393030313461653036616531353635623862613062663561323633
35343339643238633666316334636166666633336566623361323733616362643461333562643235
63303264623065306566383336633432666563373434333735656636616439313036626432323231
39306661333633656635316463366262323762303962393739643666616535633233643163303633
32393130393665613837373732666364313732323136616138396663656565666463303931326536
37346535303039316161363163356564386131333736313364386466646366363466636665383838
33623535346430333430326463376639626538356663346535613730316233633664363337373932
66643761303762613861393865323339393732373565343164383437373931313634643936376262
37613961623737656465396234313935663937623935653631323064623337336134346362336634
33323037316536653734343361653535303330326537663938663539656561333536323362386333
32633864646532383335316434613335353237633164363761613935306463633239306437656339
65373465626562333565303039376361376334366566613739656539643362353065393131336335
30323336363239313936383663353834663135643632626630616335396131333239653632663930
36613664323639616539626435346564303166623233346364396239333264373837643437356164
62613137616538396131383535313736343866366631623233626537343735626533636437376339
61623638633561386633623432396264663337643664323766663363333036626631366135363532
39663662326636393733643434393637303930666335333735363133666564343236303834323363
33323232346431353863396239613062356634303561613963393366616531366633343137653936
38616366623964386630393335646662343534633530316263373739383232383361383330336333
35623239383262616563653931363135356335303534303233363164656637343737323164336135
64316635616533633336336438353836363965343362626136623665393534363136376432333064
65366135373730623763616566623438653830303665323861363333313536376430646262613339
30373939316238623636373666323063653563306231353065303263653432366134653261623263
31653030336130363836613038366666633734353434636132373433353331323136643736393336
64383332356630393239613031653566613239373262613837366637653861633063393735396432
36336664633563393532653834386363666636316235373730653763613065383839656364353838
62623239396562643039306331306538343061333633633165383036663830656363386136613630
38313164386266633134396630613861356232393039353032353737323861653832363961356638
33623561366262383265343735366366383231336463626566636137613039653664616361396238
30356363306336313566336663306635336135653861616331323966313665633065633130336434
34636539633537323165353439613061623064316532313434396539396335636435633534663334
37323766316638653164653666316564643236353964626135373036623466646263323065656663
30626261386437333136353363316336626162656630376233366661303464643565323161353562
32303734646530396530393736366632383633333139343034383532383430616433666230346362
61636636356636663735396134386163333065636439316466376433613339383231663637626463
38613739656431346630393761646638643435376166383234306139356566653837613939336161
36613462333563656161313139373232663133373732343738653263616465316561666638666565
34396233376339353539653939666161626562666631306633316634343437336263656164383464
66313831343362323961383539386661346565643933663131633632636165643965323064343935
66333361666436656238623332353933316432323737373833333835393861613066393362313363
39313561646661336363373531613835313336646666666563333931366635636136616365303862
37303733633565363439353365626465616464386664326634346665363561653765366164336238
37623362633565313834343231663439613336653638636137383232396365373532356264343665
31653536363263376361353731653737663030383235333861303233333939343636646131333162
37303332613330386133363238356230626434643864383761653461336132616636373365313865
38396633303035343037613731663433386538646363643730656339346239363337396163663535
31373935333264666434653765306263623230666265623230306562393366376331663564653137
37333836393166386436383532663564646430333635303636383661346134623736303133616236
66353766616130333138363330386533663737626262333031376134626434336637313331306137
63366138313739663635656335326237323334643263633236353337366437643035326161303065
37383632653037353733633732356235666162653539336366323938346433623331326431616434
33363635353135663862386235396333623964616163373136343565396666393730323337623938
65373465623937643132313432636631343630353738323633356633306366636635643732333731
63336639303136316364633965336430333965306536656364613437623166643430366237306466
61333534316665333538326637306362383861356636366662613731666338303331643039343436
36646366653130383965633733633237653934616363656234373931393966323135306531333561
35656565363663326161396566353363363631356438666266393463383430653863316466373261
61316433393636653733336130616466663631646565353533383164653033643661383334623464
39343165653962366434326135626462343864353363323663386637373130643861613330343330
64656163363833356336653766386339653232336366656330326266613839323735323064333738
39613761353530346632363932356532623566353238346230383339613365343332623930366361
65623333653637366132636266333762633666653335313834316363666562663165386263373863
62396133616130653838323032663232633034326235303764393966326563313036656662643136
64626465323034346435633037633762306631393961633166313631383862636232396630333439
62353035656264653035643266616535333862656339323737333263313865663332333730356438
39303564336364393534306361333365353932616461333633303364323365356239626634316431
63626263396563313132383964376133376234633331633635306436613861313631636439313564
32393834363962636336376136346362633133323265343462393065383330393061386133663961
39616137343663636437396432356633353064373637613435666439623132306664366661383931
34383465313865373435353766636136313962373561326636626434656334646631333565366461
39353037663362666130323731343835646564393339333362646564333763366666663632326562
64346539313338323434313539323331613830623661313165373234316562313561343939666337
38343935633734326431353037303831383239656630313163626231373265346336343537336537
62376630353737663566656637626334336337343433303339616632323763633032633838346664
64306462363061396631393564303730383037613462643565353061656438343531306331313163
31383137366663636162386461623364333734613166336536643230633266656262653861386434
62633463613238333030663834353862356633333536623237326135333133386236333861643839
6339
